# deploy/docker/configs/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid

[program:postgres]
command=/usr/lib/postgresql/17/bin/postgres -D /var/lib/postgresql/data -c config_file=/etc/postgresql/17/main/postgresql.conf
user=postgres
autostart=true
autorestart=true
stdout_logfile=/var/log/postgres.log
stderr_logfile=/var/log/postgres.err.log
priority=10

[program:ocserv]
command=/usr/sbin/ocserv -c /etc/ocserv/ocserv.conf
autostart=true
autorestart=true
stdout_logfile=/var/log/ocserv.log
stderr_logfile=/var/log/ocserv.err.log
priority=20
depends_on=postgres

[program:admin-api]
command=/usr/local/bin/admin-api
autostart=true
autorestart=true
stdout_logfile=/var/log/admin-api.log
stderr_logfile=/var/log/admin-api.err.log
priority=30
depends_on=postgres

[program:customer-api]
command=/usr/local/bin/customer-api
autostart=true
autorestart=true
stdout_logfile=/var/log/customer-api.log
stderr_logfile=/var/log/customer-api.err.log
priority=40
depends_on=postgres

[program:log-parser]
command=/usr/local/bin/log-parser
autostart=true
autorestart=true
stdout_logfile=/var/log/log-parser.log
stderr_logfile=/var/log/log-parser.err.log
priority=50
depends_on=admin-api

[program:user-manager]
command=/usr/local/bin/user-manager
autostart=true
autorestart=true
stdout_logfile=/var/log/user-manager.log
stderr_logfile=/var/log/user-manager.err.log
priority=60
depends_on=admin-api

[program:telegram-bot]
command=/usr/local/bin/telegram-bot
autostart=true
autorestart=true
stdout_logfile=/var/log/telegram-bot.log
stderr_logfile=/var/log/telegram-bot.err.log
priority=70
depends_on=admin-api

[program:nginx]
command=nginx -g 'daemon off;'
autostart=true
autorestart=true
stdout_logfile=/var/log/nginx.log
stderr_logfile=/var/log/nginx.err.log
priority=80
depends_on=admin-api,customer-api
